# custom_voice_agent/supervisord.conf
; supervisord.conf — Custom Voice Agent
;
; Runs the Python backend and Next.js frontend together inside a single
; Docker container.  Both processes are restarted automatically on failure.
;
; Backend port  : 8080  (HTTP token/health endpoint)
; Frontend port : 3000  (Next.js)

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:backend]
command=python -m custom_voice_agent.main
directory=/app
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/var/log/supervisor/backend.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/backend.err
stderr_logfile_maxbytes=10MB
environment=PYTHONUNBUFFERED="1"

[program:frontend]
command=node_modules/.bin/next start -p %(ENV_PORT)s
directory=/app/frontend
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/var/log/supervisor/frontend.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/frontend.err
stderr_logfile_maxbytes=10MB
environment=PORT="3000",NODE_ENV="production"
