# custom_voice_agent/supervisord.conf
; supervisord.conf — Custom Voice Agent
;
; Runs the Python backend inside the production container.
;
; Backend port  : 8080  (health/readiness endpoint)

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:backend]
command=python -m custom_voice_agent.main
directory=/app
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/var/log/supervisor/backend.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/backend.err
stderr_logfile_maxbytes=10MB
environment=PYTHONUNBUFFERED="1",PYTHONPATH="/app/custom_voice_agent/src:/app",PORT="8080"
